# systemd/vtg-wifi-fallback.service
# vtg-wifi-fallback.service
[Unit]
Description=VTG Wifi Fallback Check

[Service]
Type=oneshot
ExecStart=/home/vtg/vtg_gen2/services/wifi_fallback.sh